# rc.local
#!/bin/sh
sysctl -p
service supervisor start
service apache2 start

