# rc.local
#!/bin/sh
sysctl -p
service supervisor start
service apache2 start
export LC_ALL="zh_CN.UTF-8"

